# archive/tools/com.grab_rebound.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grab_rebound.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mattsmacair/Desktop/Coding/量化模型/抓反弹策略/auto_update.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/grab_rebound_update.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/grab_rebound_update.log</string>
</dict>
</plist>
